# docker-flask/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:repp]
command python  /var/www/repp/run.py
